# code/explorations/codeBB/bb_source_code/cloud9/autorunBT.service
[Unit]
After=network.target

[Service]
ExecStart=/var/lib/cloud9/examples/autorunBT.sh

[Install]
WantedBy=default.target
